# generic/centos7/docrane.ini
[program:docrane]
 autorestart=true
 autostart=true
 command=/bin/docrane -v -p etcd /docrane
 redirect_stderr=true
 stdout_logfile=/var/log/docrane/docrane.log
 stdout_logfile_maxbytes=5MB
 stdout_logfile_backups=10
